# docker/alpine-apache-request-catcher/root/etc/services.d/httpd/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/httpd -k start -DFOREGROUND;
