# docker/alpine-python3/run/s6/services/python3/run
#!/usr/bin/with-contenv sh
cd /run
exec python3 /run/service.py;
