# docker/alpine-gotty-root-lscript/root/etc/services.d/gotty/run
#!/usr/bin/with-contenv sh
cd /root
exec /usr/local/bin/gotty --port 8080 --permit-write /bin/bash
